# supervisord.conf
[supervisord]
nodaemon=true
[program:sshd]
command=/usr/sbin/sshd -D

[program:mysqld]
command=/etc/init.d/mysqld start
